# LaunchDaemons/com.elyanlabs.security-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elyanlabs.security-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/elyan-monitor</string>
        <string>60</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/elyan/monitor_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/elyan/monitor_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/var/log/elyan</string>

    <key>UserName</key>
    <string>root</string>

    <key>GroupName</key>
    <string>wheel</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
